# djangofloor/templates/djangofloor/commands/supervisor.conf
# /etc/supervisor/conf.d/{{ project_name }}.conf
[program:{{ project_name }}-web]
command={{ gunicorn }} -D
user={{ username }}
autostart=true
autorestart=true
redirect_stderr=true

{% if use_celery %}
[program:{{ project_name }}-worker]
command={{ celery }} worker -D
user={{ username }}
autostart=true
autorestart=true
redirect_stderr=true
{% endif %}
