# djangofloor/templates/djangofloor/commands/upstart.conf
description "{{ project_name }}"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5

script
                NAME={{ project_name }}
                exec {{ gunicorn }} -D
{% if use_celery %}
                exec {{ celery }} worker -D
{% endif %}
end script